# .config/crontab
@reboot syncthing -no-browser &>> ${HOME}/.local/var/log/syncthing/syncthing.log
